# add/sv/phd_daemons.sv.conf
[program:phd_daemons]
#user=phabricator
command=/srv/phabricator/bin/phd restart

