# add/sv/sshd.sv.conf
[program:sshd]
command=/usr/bin/sshd -D
