# add/sv/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/usr/bin/mysqld_safe

[program:startup]
exitcodes=0
command=/bin/bash /srv/startup.sh
