# add/sv/supervisord.conf
[supervisord]
nodaemon=true

[program:startup]
exitcodes=0
command=/bin/bash /srv/startup.sh
